# corplink/root/etc/s6-overlay/s6-rc.d/corplink-service/run
#!/bin/sh
exec /opt/Corplink/corplink-service
